# etc/systemd/system/xorg@.service.d/xephyr-3seats.conf
[Service]
ExecStartPost=/usr/bin/xrandr -d :%i --output VGA --auto --right-of LVDS
