# etc/systemd/system/x-daemon-Nseats@.service.d/xephyr-5seats.conf
[Service]
ExecStartPost=/usr/bin/xrandr -d %i.0 --output VGA --auto --right-of LVDS ; /usr/bin/xrandr -d %i.1 --output VGA --auto --right-of LVDS
